# crontab
00 15 * * * bundle exec rake senryu:send_ranking