# crontab
25 19 * * * bundle exec rake senryu:send_ranking